# hoover-dryer/rootfs/etc/services.d/hoover-dryer/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Declare variables
declare MQTT_HOST
declare MQTT_PORT
declare MQTT_USERNAME
declare MQTT_PASSWORD

export MQTT_HOST=$(bashio::services mqtt "host")
export MQTT_USER=$(bashio::services mqtt "username")
export MQTT_PORT=$(bashio::services mqtt "port")
export MQTT_PASSWORD=$(bashio::services mqtt "password")

bashio::log.info "Starting Node server..."

## Run your program
exec /usr/bin/hoover-dryer
